# kiki-supervisor.conf
[program:kiki]
command=/app/venv/bin/gunicorn -b 10.0.0.30:8000 -w 4 kiki:app
directory=/app
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
